# deploy/otbasy-bot.service
[Unit]
Description=Otbasy Telegram Bot (webhook)
After=network.target

[Service]
User=otbasy
WorkingDirectory=/opt/otbasy/app
EnvironmentFile=/opt/otbasy/app/.env
ExecStart=/opt/otbasy/venv/bin/uvicorn bot:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=120
StartLimitBurst=5
MemoryMax=256M
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
